# etc/supervisord.d/vdradmin.ini
[program:vdradmin]
command=vdradmind -n
stdout_logfile = /var/log/supervisor/%(program_name)s.log
redirect_stderr=true
#stdout_syslog=true
